# infra/devpanel-worker.service
[Unit]
Description=DevPanel Agent Worker (runs on agents node, spawns claude -p)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/projects/dev-panel
ExecStart=/usr/bin/node src/worker/index.js
Restart=always
RestartSec=5

# Runtime
Environment=NODE_ENV=production
Environment=PROJECT_ROOT=/home/deploy/projects/dev-panel
Environment=DEVPANEL_STORAGE=/home/deploy/projects/dev-panel/storage
Environment=WORKER_CONCURRENCY=3
Environment=WORKER_API_PORT=3099
Environment=PATH=/home/deploy/.bun/bin:/home/deploy/.local/bin:/home/deploy/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

# Cross-node connectivity (services node exposes redis + postgres on LAN)
Environment=REDIS_HOST=10.0.0.2
Environment=REDIS_PORT=6379
Environment=PG_HOST=10.0.0.2
Environment=PG_PORT=5432
Environment=PG_USER=affine
Environment=PG_DATABASE=agent_memory
Environment=AGENT_MEMORY_NAMESPACE=dev-panel
Environment=VOYAGE_MODEL=voyage-code-3
Environment=WORKER_EVENTS_URL=https://devpanl.dev/api/admin/events/publish
# Mirror master-db writes into the services dashboard SQLite (bridge until
# orchestration tables move to shared Postgres). Worker keeps its local SQLite
# for engine reads; every createInstance/updateInstance/appendEvent/logStep
# fires a best-effort POST to DEVPANEL_REMOTE_MASTER.
# Mirror is only used while DEVPANEL_PG_ORCHESTRATION is off — once the flag
# flips, writes go straight to shared pg and this variable is ignored.
Environment=DEVPANEL_REMOTE_MASTER=https://devpanl.dev/api

# Orchestration state backend: "0" = sqlite + mirror bridge (legacy),
# "1" = shared Postgres (migration 003 tables on 10.0.0.2 agent_memory db).
# Must match the devpanel-api container's DEVPANEL_PG_ORCHESTRATION. Flip to
# "1" only after scripts/backfill-sqlite-to-pg.mjs has replayed sqlite rows
# into pg, and after services API has also been flipped (services-reads-first
# avoids a brief blind-dashboard window).
Environment=DEVPANEL_PG_ORCHESTRATION=0

# Plane REST config (non-secrets). PLANE_API_KEY is in .env.agent.
Environment=PLANE_BASE_URL=https://plane.devpanl.dev
Environment=PLANE_WORKSPACE_SLUG=devpanl
Environment=PLANE_PROJECT_ID=d2522fed-e3f2-4eeb-9077-6445261752c1

# Backlog puller defaults — flip BACKLOG_PULL_ENABLED to "true" to start the
# continuous Plane-Todo → workflow loop. LABEL gate keeps only items you've
# explicitly tagged "agent-ready" in-scope for the agent team.
Environment=BACKLOG_PULL_ENABLED=true
Environment=BACKLOG_PULL_LABEL=agent-ready
Environment=BACKLOG_PULL_STATE_GROUPS=unstarted
Environment=BACKLOG_PULL_MAX_PER_TICK=3
Environment=BACKLOG_PULL_INTERVAL_MS=900000

# Secrets come from a hand-managed file on the agents node (not CI-deployed):
#   /home/deploy/projects/dev-panel/.env.agent
# Required keys: PG_PASSWORD, VOYAGE_API_KEY, ADMIN_API_KEY, GITHUB_TOKEN,
#                TELEGRAM_BOT_TOKEN, TELEGRAM_CHAT_ID, PLANE_API_KEY
EnvironmentFile=-/home/deploy/projects/dev-panel/.env.agent

StandardOutput=journal
StandardError=journal
SyslogIdentifier=devpanel-worker

[Install]
WantedBy=multi-user.target
