# infra/shelly.service
[Unit]
Description=Shelly — Claude Code + telegram plugin (tmux-backed)
Documentation=file:///home/deploy/projects/dev-panel/CLAUDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/projects/dev-panel
Environment=HOME=/home/deploy
Environment=PATH=/home/deploy/.bun/bin:/home/deploy/.local/bin:/home/deploy/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
# devpanel-mcp (spawned by Shelly's claude session) connects to BullMQ Redis
# and pgvector Postgres on the services VPS via the private LAN. Without
# these, src/mcp/server.js falls back to the public IP 77.42.46.87 which is
# firewalled — every list_jobs/enqueue_job/memory_search call hangs forever
# and Shelly looks "stuck". Mirror what devpanel-worker.service sets.
Environment=REDIS_HOST=10.0.0.2
Environment=REDIS_PORT=6379
Environment=PG_HOST=10.0.0.2
Environment=PG_PORT=5432
Environment=DEVPANEL_STORAGE=/home/deploy/projects/dev-panel/storage
Environment=PROJECT_ROOT=/home/deploy/projects/dev-panel
EnvironmentFile=-/home/deploy/projects/dev-panel/.env.agent

# Make sure log destination exists.
ExecStartPre=-/bin/mkdir -p /home/deploy/logs
ExecStartPre=-/bin/chown deploy:deploy /home/deploy/logs

# Always start from a clean tmux session.
ExecStartPre=-/usr/bin/tmux -L deploy kill-session -t shelly

# Inbound Telegram is delivered by `telegram-multi` (stdio MCP server + grammy
# long-poller, wired in ~/.mcp.json under server name `telegram`). `--channels`
# uses the `server:<name>` form for manually configured MCP servers (vs.
# `plugin:<name>@<marketplace>` for plugin-marketplace channels). Previously
# this used `plugin:telegram@claude-plugins-official`, which spawned a SECOND
# poller fighting telegram-multi for the same bot token (Telegram returns 409
# Conflict when two pollers share a token, both eventually go silent). Switched
# to `server:telegram` 2026-04-27 so telegram-multi is the sole poller.
# `--dangerously-load-development-channels` is required by Claude when --channels
# uses the `server:<name>` form (channels backed by manually-configured MCP
# servers rather than vetted plugin-marketplace plugins). Without it Claude
# refuses to subscribe and Shelly receives nothing — even though telegram-multi
# is happily long-polling.
ExecStart=/usr/bin/tmux -L deploy new-session -d -s shelly "/usr/bin/claude --channels server:telegram --dangerously-load-development-channels --dangerously-skip-permissions"

# 1) Dismiss claude's first-run "trust this folder?" prompt.
# 2) Mirror the pane to the systemd journal under tag "shelly" so journalctl
#    -t shelly survives a session kill (tmux pty would otherwise lose the
#    crash trace). pipe-pane appends to a fresh stream on every restart.
# 3) Also tee to a plain file so a frustrated human doesn't need journalctl
#    to read the last 200 lines.
ExecStartPost=/bin/sh -c "sleep 6 && \
  /usr/bin/tmux -L deploy send-keys -t shelly Enter && \
  /usr/bin/tmux -L deploy pipe-pane -O -t shelly:0.0 'tee -a /home/deploy/logs/shelly.log | systemd-cat -t shelly -p info'"

ExecStop=-/usr/bin/tmux -L deploy kill-session -t shelly

[Install]
WantedBy=multi-user.target
